# .config/systemd/user/package-list.service
[Unit]
Description=Get packages installed on Arch system

[Service]
Type=oneshot
ExecStart=/home/luqman/bin/get-installed-arch
