# .config/systemd/user/polybar.target
[Unit]
Description=Polybar
Documentation=man:polybar(1)
Wants=polytop.service polybot.service
